# etc/profile.d/notify_log.sh
# Notify user loggin
DISPLAY=:0.0 notify-send -t 4000 "User logged in " `whoami`
